# .config/systemd/user/compton.service
[Unit]
Description=Compositing manager

[Service]
Environment="DISPLAY=:0"
ExecStart=/usr/bin/compton

[Install]
WantedBy=xsession.target
