# .config/systemd/user/keymap.target
[Unit]
Description=Keymap for Laptop
Requires=laptop-keymap.service xcape.service
